# logsentry.cron
00 * * * *	root	/bin/sh /usr/sbin/logcheck.sh
